# deployment/systemd/workmind.service
# Workmind 后端服务（Daphne ASGI，HTTP + WebSocket）
# 安装：cp /data/workmind/deployment/systemd/workmind.service /etc/systemd/system/
#       systemctl daemon-reload && systemctl enable workmind && systemctl start workmind
# 以 root 运行，省去 ADB/目录权限配置（生产建议改为 nginx 或专用用户）

[Unit]
Description=Workmind Daphne ASGI Server
After=network.target mysql.service redis.service mariadb.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/data/workmind/backend
Environment="PATH=/data/workmind/backend/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/data/workmind/.env
ExecStart=/data/workmind/backend/venv/bin/daphne -b 0.0.0.0 -p 8009 backend.asgi:application
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
